# sched/nightly-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler unit for NIGHTLY.

  Two ways to install:

  A) Import this file via Task Scheduler GUI:
       Task Scheduler -> Action -> Import Task... -> select this file.

  B) Import via PowerShell (run as admin):
       Register-ScheduledTask -Xml (Get-Content -Path 'nightly-task.xml' -Raw) -TaskName 'NIGHTLY'

  C) Or build the task inline without this XML:
       $Action  = New-ScheduledTaskAction -Execute 'claude' -Argument '-p ''/nightly'''
       $Trigger = New-ScheduledTaskTrigger -Daily -At '10:00PM'
       Register-ScheduledTask -TaskName 'NIGHTLY' -Action $Action -Trigger $Trigger -RunLevel Highest

  Adjust the Command path if `claude` isn't on PATH for the SYSTEM/USER
  account running the task. Easiest way: `(Get-Command claude).Source` in
  PowerShell to find the absolute path.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>NIGHTLY autoresearch loop - improves Claude Code substrate overnight via personal session-history benchmark.</Description>
    <URI>\NIGHTLY</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-01T22:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT2H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c claude -p "/nightly" &gt;&gt; "%USERPROFILE%\.claude\nightly\logs\cron.log" 2&gt;&amp;1</Arguments>
    </Exec>
  </Actions>
</Task>
